# roles/irods3-docker/templates/profile.d/irods.sh.j2
export IRODS_HOME=/opt/irods-legacy-3.3.1

# modify path to include iRODS server and icommands
export PATH="$IRODS_HOME/iRODS/server/bin/cmd:$IRODS_HOME/iRODS/clients/icommands/bin:$PATH"
